# httpd.service
.include /lib/systemd/system/httpd.service
[Service]
PIDFile=/var/run/httpd/httpd.pid

