# httpd.service
[Service]
PIDFile=/var/run/httpd/httpd.pid
