# rhel-domainname.service
[Service]
ExecStart=/bin/true
Type=oneshot
RemainAfterExit=yes
